# redis.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=99
USE_PROCD=1
start_service() {	
        procd_open_instance "redis-server"
	procd_set_param command "/bin/redis-server"
	procd_append_param command  "/etc/redis.conf"
	procd_set_param respawn
	#procd_set_param file "$config_file"	 
	procd_close_instance
	 
 
}

